# check-updates.service
[Unit]
Description=Vérifie les mises à jour apt et envoie sur Discord
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/check-updates.sh
